# scripts/com.mahesh.driftai-jobs.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.mahesh.driftai-jobs</string>

  <!-- Fetch jobs, then commit+push if changed (triggers Netlify redeploy) -->
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/maheshkunasani/IdeaProjects/resume-marketing-site/scripts/daily-update.sh</string>
  </array>

  <!-- Ensure node + git are on PATH when launchd runs this -->
  <key>WorkingDirectory</key>
  <string>/Users/maheshkunasani/IdeaProjects/resume-marketing-site</string>

  <!-- Only PATH here. Secret keys live in scripts/.env (gitignored),
       which daily-update.sh loads at runtime — so no keys touch GitHub. -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/maheshkunasani/.nvm/versions/node/v24.16.0/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <!-- Run twice every day: 6:00 AM and 6:00 PM -->
  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Hour</key><integer>6</integer>
      <key>Minute</key><integer>0</integer>
    </dict>
    <dict>
      <key>Hour</key><integer>18</integer>
      <key>Minute</key><integer>0</integer>
    </dict>
  </array>

  <!-- Log into the project dir, same as stock-digest (launchd.log is gitignored) -->
  <key>StandardOutPath</key>
  <string>/Users/maheshkunasani/IdeaProjects/resume-marketing-site/launchd.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/maheshkunasani/IdeaProjects/resume-marketing-site/launchd.log</string>
</dict>
</plist>
